# install/craftpanel.service
# Reference unit file. The installer (install.sh) writes this automatically,
# adjusted to the chosen port and repo.
[Unit]
Description=ComputeBox Craftpanel (Minecraft server panel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=craftpanel
Group=craftpanel
ExecStart=/var/lib/craftpanel/bin/craftpanel -data /var/lib/craftpanel -addr :8420
Restart=always
RestartSec=3

# Only ever signal the panel itself. On systemctl stop it shuts its Minecraft
# servers down gracefully; during self updates it exits without touching them
# and the new panel instance reattaches to the running servers.
KillMode=process
TimeoutStopSec=90

# Hardening. MemoryDenyWriteExecute stays off because the Java JIT needs it.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# The binary lives inside the data directory so the panel can self update.
ReadWritePaths=/var/lib/craftpanel
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
